# pkg/scripts/postinstall
#!/bin/sh
set -e
/bin/launchctl load /Library/LaunchDaemons/ch.roe.xnumon.plist
